# verbose.service
[Unit]
Description=Verbose - Voice-to-text daemon
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/nax/tools/verbose/verbose.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
