# verbose.service
[Unit]
Description=Verbose - Voice-to-text daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/tools/verbose/verbose.py
Restart=on-failure
RestartSec=5
Environment="DISPLAY=:0"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
